# provisioning/files/init-uwsgi.conf
description "uwsgi tiny instance"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --master --die-on-term --emperor /etc/uwsgi --emperor-pidfile /var/run/uwsgi.pid
